# scripts/systemd/novnc.service
# Servicio systemd para noVNC (proxy web hacia x11vnc)
# Copiar a /etc/systemd/system/ y ajustar User y ruta de noVNC según el servidor.
#
# Uso:
#   sudo cp scripts/systemd/novnc.service /etc/systemd/system/
#   sudo nano /etc/systemd/system/novnc.service   # ajustar User y NOVNC_PATH
#   sudo systemctl daemon-reload
#   sudo systemctl enable novnc
#   sudo systemctl start novnc
#
# Conviene arrancar después de mop-dga para que x11vnc esté ya escuchando en 5900.

[Unit]
Description=noVNC proxy (acceso web al escritorio virtual MOP DGA)
After=network-online.target mop-dga.service
Wants=network-online.target

[Service]
Type=simple
User=datos
Group=datos
# Ajustar al camino real de noVNC en el servidor (ej. /home/datos/noVNC)
WorkingDirectory=/home/datos/noVNC
ExecStart=/home/datos/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 6080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
